# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; ===================
; PHP-FPM
; ===================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; ===================
; Nginx
; ===================
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; ===================
; Inertia SSR (Node). If it fails, Laravel falls back to client render
; (INERTIA_SSR_THROW_ON_ERROR=false), so the site stays up regardless.
; ===================
[program:ssr]
command=/bin/sh -c "while [ ! -f /var/www/html/bootstrap/cache/config.php ]; do sleep 1; done && /usr/local/bin/php /var/www/html/artisan inertia:start-ssr"
user=www-data
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ===================
; Queue worker (database driver). Handles async jobs: order notifications.
; WAL + busy_timeout on SQLite keeps it happy
; alongside web writes. Move to Redis + Horizon if/when throughput demands.
; ===================
[program:worker]
command=/bin/sh -c "while [ ! -f /var/www/html/bootstrap/cache/config.php ]; do sleep 1; done && /usr/local/bin/php /var/www/html/artisan queue:work --tries=3 --sleep=3 --max-time=3600 --rest=1"
user=www-data
autostart=true
autorestart=true
priority=40
stopwaitsecs=70
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; ===================
; Scheduler. schedule:work runs due tasks (abandoned-payment sweep, etc)
; without needing cron inside the container; it inherits env correctly and
; logs to stdout like everything else.
; ===================
[program:scheduler]
command=/bin/sh -c "while [ ! -f /var/www/html/bootstrap/cache/config.php ]; do sleep 1; done && /usr/local/bin/php /var/www/html/artisan schedule:work"
user=www-data
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
